# cccc.conf
[program:cccc]
command=python3 /home/curator2/cccc/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/cccc/test.err.log
stdout_logfile=/home/curator2/cccc/test.out.log